# packages/cc-gateway/systemd/duraclaw-cc-gateway.service
[Unit]
Description=Duraclaw CC Gateway - Claude Code API Gateway
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/data/projects/duraclaw
ExecStart=/home/ubuntu/.bun/bin/bun run packages/cc-gateway/src/server.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/data/projects/duraclaw/.env
Environment=NODE_ENV=production
Environment=CC_GATEWAY_PORT=9877
Environment=CLAUDE_CODE_ENABLE_TASKS=true
Environment=PATH=/home/ubuntu/.local/bin:/home/ubuntu/.bun/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/ubuntu

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/data/projects /home/ubuntu/.cc-gateway /home/ubuntu/.claude /tmp

[Install]
WantedBy=multi-user.target
